# airootfs/etc/systemd/system/optimize-rt-blankscreen.service
[Unit]
Description=Real Time Configuration

[Service]
# Disable screen blanking
ExecStart=/usr/bin/xset -dpms; xset s off

[Install]
WantedBy=multi-user.target 